# misc/dynaswap.service
[Unit]
Description=Dynamic Swapping under load

[Service]
ExecStart=/bin/dynaswap -c /etc/dynaswap.conf

[Install]
WantedBy=multi-user.target